# mem0-api.service
[Unit]
Description=mem0 injection proxy with quality filtering
After=network.target qdrant.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt
ExecStart=/usr/bin/python3 /opt/mem0_proxy.py
Restart=on-failure
RestartSec=5
Environment=MEM0_PROXY_HOST=0.0.0.0
Environment=MEM0_PROXY_PORT=8765
Environment=MEM0_CONFIG_PATH=/root/.mem0/config.json
Environment=MEM0_PROXY_LOG_LEVEL=INFO
# Optional: Custom quality filter prompts (uncomment to use)
# Environment=MEM0_PROXY_FACT_EXTRACTION_PROMPT=/etc/mem0/fact_extraction.txt
# Environment=MEM0_PROXY_UPDATE_MEMORY_PROMPT=/etc/mem0/update_memory.txt

[Install]
WantedBy=multi-user.target
